# setup_wifi/static_wpa_network.service.template
[Unit]
Description=Network Connectivity
Wants=network.target
Before=network.target
Wants=${network_device}
Before=${network_device}
BindsTo=sys-subsystem-net-devices-${interface}.device
After=sys-subsystem-net-devices-${interface}.device

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=${if_destination}
ExecStart=/sbin/ip link set dev ${interface} up
ExecStart=/usr/sbin/wpa_supplicant -B -i ${interface} -c ${wpa_destination}
ExecStart=/usr/sbin/ip addr add ${address}/${netmask} broadcast ${broadcast} dev ${interface}
ExecStart=/usr/sbin/ip route add default via ${gateway}
ExecStop=/usr/sbin/ip addr flush dev ${interface}
ExecStop=/usr/sbin/ip link set dev ${interface} down

[Install]
WantedBy=multi-user.target
